# deploy/relay.service
# systemd service — onchato relay
# Skopiuj do: /etc/systemd/system/onchato-relay.service
# Aktywuj:
#   sudo systemctl daemon-reload
#   sudo systemctl enable onchato-relay
#   sudo systemctl start onchato-relay
# Logi:
#   sudo journalctl -u onchato-relay -f

[Unit]
Description=onchato libp2p relay
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/onchato
ExecStart=/usr/bin/node src/relay.mjs --port 9001 --pass ZMIEN_MNIE --host onchato.com
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
